# root/etc/services.d/Mysql/run
#!/usr/bin/with-contenv bash

exec /app/Mysql/mysql run -confdir /config/Mysql